# config/passmaster-shutdown.conf
# Passmaster Upstart Configuration - Shutdown

description "Shutdown Passmaster"

start on (runlevel [016])

exec su - ryan -c 'cd /home/ryan/passmaster && script/wrapper script/unicorn stop'
